# deployment/config-files/s6-rc.d/oneshot-artisan-migrate/up
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

php /app/artisan migrate --force