# deployment/config-files/s6-rc.d/longrun-nginx/run
#!/command/execlineb -P
####################################
# start up service
####################################

# keep original environment
with-contenv

# start service in foreground
nginx -g "daemon off;" -c /etc/nginx/nginx.conf.replaced