# deployment/config-files/s6-rc.d/longrun-artisan-queue-work/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

# start service in foreground
/app/artisan queue:work