# deployment/config-files/s6-rc.d/oneshot-composer-install/up
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

cd /app
foreground {
  if -n { test -d /var/www/.composer }
  mkdir -p /var/www/.composer
}
foreground {
  composer config --global cache-dir /composer_cache
}
foreground {
  composer config --global preferred-install.sandstorm/* source
}
composer install