# deployment/config-files/s6-rc.d/longrun-artisan-reverb-start/run
#!/command/execlineb -P

####################################
# start up service
####################################

# keep original environment
with-contenv

# start service in foreground
# for debugging, add "--debug"
/app/artisan reverb:start --port=8899
